# systemd/borg-backup.timer
[Unit]
Description=Create backup using Borg Backup

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
